# shadowsocks-client/etc/supervisor/shadowsocks.conf
[program:shadowsocks]
command=sslocal -c /etc/shadowsocks.json
autorestart=true
user=nobody
