# shadowsocks-client/etc/supervisor/shadowsocks.conf
[program:shadowsocks]
command=ssclient -c /etc/shadowsocks.json
autorestart=true
user=nobody
